# supervisord.conf
[program:certbotbot]
command=/certbotbot.sh
startsecs=0
autorestart=unexpected
